# manjaro/thinkfan.service
[Unit]
Description=Thinkfan

[Service]
ExecStart=/usr/sbin/thinkfan -n

[Install]
WantedBy=multi-user.target
