# scripts/backup/backup.timer
[Unit]
Description=Run backup to SD daily

[Timer]
OnCalendar=daily
Persistent=true
 
[Install]
WantedBy=timers.target
