# scripts/backup/backup.service
[Unit]
Description=Backup home to SD

[Service]
Type=simple
ExecStart=/home/gcg/dotfiles/scripts/backup/daily.sh
